# deploy/codex-quota-api.service
[Unit]
Description=ZongTech Codex weekly quota API
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
DynamicUser=yes
WorkingDirectory=/opt/codex-quota-api
Environment=NODE_ENV=production
EnvironmentFile=-/etc/codex-quota-api.env
ExecStart=/usr/bin/node /opt/codex-quota-api/server/codex-quota-server.js
Restart=on-failure
RestartSec=3s
TimeoutStopSec=15s
UMask=0077

NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
CapabilityBoundingSet=
AmbientCapabilities=
LockPersonality=true
RemoveIPC=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
